# eye-reminder.service
[Unit]
Description=Eye Reminder Service

[Service]
ExecStart=/etc/eye-timer/notification-script.sh
